# templates/upstart-elasticsearch.conf.j2
#!upstart

description "Start elasticsearch"

start on runlevel [2345]
stop on runlevel [!2345]

console log

setuid socorro

exec "{{ ansible_env.HOME }}/elasticsearch-1.4.5/bin/elasticsearch"
